# deploy/com.udhay.eduwallet-autodeploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Reference copy. Installed at ~/Library/LaunchAgents/com.udhay.eduwallet-autodeploy.plist
     on jarvis. Runs the poller from the DEPLOY dir copy (refreshed by every deploy's rsync).
     Manage: launchctl bootstrap gui/501 <plist> | kickstart -k gui/501/com.udhay.eduwallet-autodeploy -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.udhay.eduwallet-autodeploy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jarvis/eduwallet/deploy/auto-deploy-poll.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/jarvis/Library/Logs/eduwallet-autodeploy.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jarvis/Library/Logs/eduwallet-autodeploy.err.log</string>
</dict>
</plist>
